# ansible/roles/gotomation/templates/gotomation.service.j2
[Unit]
Description=Gotomation
After=network.target

[Service]
ExecStart=/usr/local/bin/gotomation-{{ gotomation.arch }}-{{ gotomation.version }} -c {{ gotomation.config }}
WorkingDirectory=/etc
StandardOutput=inherit
StandardError=inherit
Restart=always
User={{ gotomation.user }}

[Install]
WantedBy=multi-user.target
